# ops/mac-mini/launchd/com.agentpay.agency-os.evolve.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.agentpay.agency-os.evolve</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/baskar_viji/Agentpay/ops/mac-mini/bin/agency-os-evolve</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/baskar_viji/Agentpay</string>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>StandardOutPath</key>
    <string>/Users/baskar_viji/.openclaw/logs/agency-os-evolve.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/baskar_viji/.openclaw/logs/agency-os-evolve.err.log</string>
  </dict>
</plist>
